# ansible/roles/hyperdx/templates/hyperdx-app.service.j2
[Unit]
Description=HyperDX Frontend (Next.js)
After=forge-metal-firewall.target network.target hyperdx-api.service
Requires=forge-metal-firewall.target

[Service]
Type=simple
User={{ ansible_user }}
WorkingDirectory={{ hyperdx_install_dir }}/packages/app
EnvironmentFile={{ hyperdx_env_file }}
EnvironmentFile=/etc/credstore/hyperdx/secrets.env
ExecStart={{ fm_bin }}/npx next start -p {{ services.hyperdx.app_port }}
Restart=on-failure
RestartSec=5
Environment=PATH={{ fm_bin }}:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
