# src/platform/ansible/roles/nftables/templates/verself-firewall.target.j2
# verself-firewall.target — abstraction point for host firewall readiness.
#
# WHY THIS EXISTS:
# On boot, systemd starts services by dependency graph — not Ansible role order.
# If a service binds its listen port before nftables loads the "policy drop"
# ruleset, that port is briefly reachable under the kernel's default ACCEPT
# policy — exposed to the internet with no filtering.
#
# Every network-facing verself service declares:
#   After=verself-firewall.target
#   Requires=verself-firewall.target
# so it cannot start until the drop policy is active.
#
# This target wraps the concrete firewall implementation (currently
# nftables.service). If the firewall backend changes, only this file
# needs updating — not every service unit.

[Unit]
Description=verself host firewall ready
Requires=nftables.service
After=nftables.service

[Install]
WantedBy=multi-user.target
