# src/platform/ansible/roles/otelcol/templates/otelcol.service.j2
[Unit]
Description=OpenTelemetry Collector
After=forge-metal-firewall.target network.target clickhouse-server.service spire-agent.service otelcol-clickhouse-spiffe-helper.service
Requires=forge-metal-firewall.target
Wants=clickhouse-server.service spire-agent.service otelcol-clickhouse-spiffe-helper.service

[Service]
Type=simple
User={{ otelcol_system_user }}
Group={{ otelcol_system_group }}
SupplementaryGroups={{ ([spire_workload_group] + otelcol_log_reader_groups) | join(' ') }}
ExecStartPre={{ fm_bin }}/spiffe-helper -daemon-mode=false -config {{ otelcol_clickhouse_helper_config_path }}
ExecStart={{ fm_bin }}/otelcol-contrib --config /etc/otelcol/config.yaml
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/otelcol/storage {{ otelcol_clickhouse_spiffe_dir }}

[Install]
WantedBy=multi-user.target
