# src/platform/ansible/roles/verself_web/templates/verself-web.service.j2
[Unit]
Description=Verself Web (Console + Docs + Policy, TanStack Start)
After=verself-firewall.target network.target sandbox-rental-service.service identity-service.service notifications-service.service electric.service electric-notifications.service
Requires=verself-firewall.target
Wants=sandbox-rental-service.service identity-service.service notifications-service.service electric.service electric-notifications.service

[Service]
Type=simple
User={{ verself_web_runtime_user }}
Group={{ verself_web_runtime_group }}
WorkingDirectory={{ verself_web_app_dir }}
EnvironmentFile=/etc/verself-web/env
ExecStart={{ verself_bin }}/node --import ./instrumentation.mts .output/server/index.mjs
Restart=on-failure
RestartSec=5
Environment=PATH={{ verself_bin }}:/usr/bin:/bin
Environment=HOME={{ verself_web_state_dir }}
CapabilityBoundingSet=
ProtectHome=true
ProtectSystem=strict
ReadWritePaths={{ verself_web_state_dir }}
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
LockPersonality=true
NoNewPrivileges=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=multi-user.target
