# src/platform/ansible/roles/rent_a_sandbox/templates/rent-a-sandbox.service.j2
[Unit]
Description=Rent-a-Sandbox Frontend (TanStack Start)
After=forge-metal-firewall.target network.target sandbox-rental-service.service electric.service
Requires=forge-metal-firewall.target
Wants=sandbox-rental-service.service electric.service

[Service]
Type=simple
User={{ ansible_user }}
WorkingDirectory={{ rent_a_sandbox_install_dir }}
EnvironmentFile=/etc/rent-a-sandbox/env
ExecStart={{ fm_bin }}/node .output/server/index.mjs
Restart=on-failure
RestartSec=5
Environment=PATH={{ fm_bin }}:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
