# ansible/roles/firecracker/templates/forge-metal-firecracker-network.service.j2
[Unit]
Description=Configure forge-metal Firecracker host networking
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/forge-metal-firecracker-network-setup
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
