# src/platform/ansible/roles/electric_letters/templates/electric-letters.service.j2
[Unit]
Description=Electric Sync Service (Letters)
After=forge-metal-firewall.target network.target postgresql.service
Requires=forge-metal-firewall.target
Wants=postgresql.service

[Service]
Type=simple
ExecStartPre=/usr/bin/podman pull {{ electric_letters_image }}
ExecStart=/usr/bin/podman run --rm --replace \
  --name electric-letters \
  --network=host \
  -e DATABASE_URL=postgresql://electric_letters:{{ electric_letters_pg_password }}@{{ services.postgresql.host }}:{{ services.postgresql.port }}/letters \
  -e ELECTRIC_PORT={{ services.electric_letters.port }} \
  -e ELECTRIC_SECRET={{ electric_letters_api_secret }} \
  -e ELECTRIC_STORAGE_DIR=/data \
  -e ELECTRIC_INSTANCE_ID=electric-letters \
  -e ELECTRIC_REPLICATION_STREAM_ID=letters \
  -e RELEASE_NAME=electric_letters \
  -v {{ electric_letters_storage_dir }}:/data \
  {{ electric_letters_image }}
ExecStop=/usr/bin/podman stop electric-letters
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
