# src/platform/ansible/roles/company/templates/company.service.j2
[Unit]
Description=Guardian Intelligence Company Site (TanStack Start)
After=verself-firewall.target network.target
Requires=verself-firewall.target

[Service]
Type=simple
User={{ company_runtime_user }}
Group={{ company_runtime_group }}
WorkingDirectory={{ company_current_link }}
EnvironmentFile=/etc/company/env
ExecStart={{ verself_bin }}/node ./app-entry.mjs
Restart=on-failure
RestartSec=5
Environment=PATH={{ verself_bin }}:/usr/bin:/bin
Environment=HOME={{ company_state_dir }}
CapabilityBoundingSet=
ProtectHome=true
ProtectSystem=strict
ReadWritePaths={{ company_state_dir }}
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
LockPersonality=true
NoNewPrivileges=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=multi-user.target
