# src/components/pomerium/templates/pomerium.service.j2
[Unit]
Description=Pomerium operator access gateway
After=verself-firewall.target network.target postgresql.service zitadel.service otelcol.service grafana.service
Requires=verself-firewall.target postgresql.service zitadel.service
Wants=otelcol.service grafana.service

[Service]
Type=simple
User=pomerium
Group=pomerium
WorkingDirectory={{ pomerium_data_dir }}
Environment=OTEL_TRACES_EXPORTER=otlp
Environment=OTEL_EXPORTER_OTLP_ENDPOINT=http://{{ pomerium_otelcol_otlp_grpc_address }}
Environment=OTEL_EXPORTER_OTLP_PROTOCOL=grpc
Environment=OTEL_RESOURCE_ATTRIBUTES=service.name=pomerium,service.namespace=verself-substrate
ExecStart={{ verself_bin }}/pomerium --config {{ pomerium_config_path }}
Restart=on-failure
RestartSec=5
LimitNOFILE=65536
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths={{ pomerium_data_dir }} {{ pomerium_log_dir }}

[Install]
WantedBy=multi-user.target
