# ansible/roles/agent/templates/forge-metal-agent.service.j2
[Unit]
Description=forge-metal CI agent
After=network.target containerd.service
Requires=containerd.service

[Service]
Type=simple
ExecStart={{ fm_bin }}/forge-metal agent run --config /etc/forge-metal/config.toml
Restart=on-failure
RestartSec=5
Environment=PATH={{ fm_bin }}:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
