# src/platform/ansible/roles/temporal/templates/temporal-server.service.j2
{% set namespace_roles_env -%}
{%- for binding in temporal_namespace_role_bindings -%}
{{ binding.spiffe_id }}|{{ binding.namespace }}|{{ binding.role }}{% if not loop.last %},{% endif %}
{%- endfor -%}
{%- endset %}
[Unit]
Description=Temporal Server
After=forge-metal-firewall.target network.target postgresql.service spire-agent.service otelcol.service
Requires=forge-metal-firewall.target
Wants=postgresql.service spire-agent.service otelcol.service

[Service]
Type=simple
User={{ temporal_system_user }}
Group={{ temporal_system_group }}
SupplementaryGroups={{ spire_workload_group }}
Environment=SPIFFE_ENDPOINT_SOCKET=unix://{{ spire_agent_socket_path }}
Environment=OTEL_SERVICE_NAME=temporal-server
Environment=OTEL_TRACES_EXPORTER=otlp
Environment=OTEL_EXPORTER_OTLP_ENDPOINT=http://{{ services.otelcol.host }}:{{ services.otelcol.grpc_port }}
Environment=OTEL_EXPORTER_OTLP_INSECURE=true
Environment=OTEL_EXPORTER_OTLP_TRACES_PROTOCOL=grpc
Environment=FM_TEMPORAL_CONFIG_PATH={{ temporal_config_path }}
Environment=FM_TEMPORAL_SERVER_SPIFFE_ID={{ spire_temporal_server_id }}
Environment=FM_TEMPORAL_FRONTEND_CLIENT_IDS={{ temporal_frontend_client_spiffe_ids | join(',') }}
Environment=FM_TEMPORAL_SYSTEM_ADMIN_IDS={{ spire_temporal_server_id }}
Environment=FM_TEMPORAL_SYSTEM_READER_IDS={{ spire_temporal_web_id }}
Environment=FM_TEMPORAL_NAMESPACE_ROLES={{ namespace_roles_env }}
WorkingDirectory={{ temporal_state_dir }}
ExecStart={{ fm_bin }}/forge-temporal-server
Restart=on-failure
RestartSec=5
CapabilityBoundingSet=
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
LockPersonality=true
ReadWritePaths={{ temporal_state_dir }}
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=multi-user.target
