# src/platform/ansible/roles/temporal/templates/temporal-proof-worker.service.j2
[Unit]
Description=Temporal Proof Worker
After=network.target temporal-server.service governance-service.service spire-agent.service otelcol.service
Wants=temporal-server.service governance-service.service spire-agent.service otelcol.service

[Service]
Type=simple
User={{ temporal_proof_user }}
Group={{ temporal_proof_group }}
SupplementaryGroups={{ spire_workload_group }}
Environment=SPIFFE_ENDPOINT_SOCKET=unix://{{ spire_agent_socket_path }}
Environment=OTEL_EXPORTER_OTLP_ENDPOINT=http://{{ services.otelcol.host }}:{{ services.otelcol.grpc_port }}
Environment=FM_TEMPORAL_FRONTEND_ADDRESS={{ services.temporal.host }}:{{ services.temporal.grpc_port }}
Environment=FM_TEMPORAL_PROOF_NAMESPACE={{ temporal_proof_namespace }}
Environment=FM_TEMPORAL_SERVER_SPIFFE_ID={{ spire_temporal_server_id }}
Environment=FM_TEMPORAL_PROOF_GOVERNANCE_URL={{ temporal_proof_governance_url }}
Environment=FM_TEMPORAL_PROOF_GOVERNANCE_SPIFFE_ID={{ spire_governance_service_id }}
Environment=FM_TEMPORAL_PROOF_SERVICE_VERSION={{ forge_metal_version }}
WorkingDirectory={{ temporal_proof_state_dir }}
ExecStart={{ fm_bin }}/temporal-proof worker
Restart=on-failure
RestartSec=5
CapabilityBoundingSet=
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
LockPersonality=true
ReadWritePaths={{ temporal_proof_state_dir }}
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=multi-user.target
