# src/platform/ansible/roles/temporal/templates/temporal-web.service.j2
[Unit]
Description=Temporal Web
After=forge-metal-firewall.target network.target temporal-server.service zitadel.service caddy.service spire-agent.service otelcol.service
Requires=forge-metal-firewall.target
Wants=temporal-server.service zitadel.service caddy.service spire-agent.service otelcol.service

[Service]
Type=simple
User={{ temporal_web_user }}
Group={{ temporal_web_group }}
SupplementaryGroups={{ spire_workload_group }}
Environment=SPIFFE_ENDPOINT_SOCKET=unix://{{ spire_agent_socket_path }}
Environment=OTEL_SERVICE_NAME=temporal-web
Environment=OTEL_TRACES_EXPORTER=otlp
Environment=OTEL_EXPORTER_OTLP_ENDPOINT=http://{{ services.otelcol.host }}:{{ services.otelcol.grpc_port }}
Environment=OTEL_EXPORTER_OTLP_INSECURE=true
Environment=OTEL_EXPORTER_OTLP_TRACES_PROTOCOL=grpc
Environment=FM_TEMPORAL_WEB_CONFIG_DIR={{ temporal_web_config_dir }}
Environment=FM_TEMPORAL_WEB_CONFIG_ENV=production
Environment=FM_TEMPORAL_FRONTEND_ADDRESS={{ services.temporal.host }}:{{ services.temporal.grpc_port }}
WorkingDirectory={{ temporal_web_state_dir }}
ExecStart={{ fm_bin }}/forge-temporal-web
Restart=on-failure
RestartSec=5
CapabilityBoundingSet=
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
LockPersonality=true
ReadWritePaths={{ temporal_web_state_dir }}
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=multi-user.target
