# src/platform/ansible/roles/firecracker/templates/homestead-smelter-host.service.j2
[Unit]
Description=Homestead Smelter Firecracker Host Agent
After=forge-metal-firewall.target network-online.target forge-metal-firecracker-network.service
Requires=forge-metal-firewall.target
Wants=network-online.target

[Service]
Type=simple
ExecStart={{ firecracker_smelter_host_binary_path }} serve --listen-uds {{ firecracker_smelter_control_socket }} --ops-uds {{ firecracker_smelter_ops_socket }} --jailer-root {{ firecracker_jailer_root }}/firecracker
ExecStartPost=/usr/bin/chgrp {{ firecracker_smelter_ops_group }} {{ firecracker_smelter_ops_socket }}
Restart=always
RestartSec=2
RuntimeDirectory=homestead-smelter
Environment=PATH=/usr/local/bin:/usr/bin:/bin
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
