# src/host-configuration/ansible/playbooks/wipe-pg-db.yml
---
# Wipe one managed PostgreSQL service database and recreate substrate state.
# Application migrations run through nomad-deploy-all during the follow-up
# Nomad deploy.
#
# Usage:
#   verself-deploy ansible run --site=prod \
#     --playbook=playbooks/wipe-pg-db.yml \
#     --inventory=$(pwd)/src/host-configuration/ansible/inventory/prod.ini \
#     --ansible-arg=-e --ansible-arg=wipe_pg_db_name=sandbox_rental
#
# Supported DBs:
# - billing: billing ledger reset, including TigerBeetle
# - sandbox_rental
# - mailbox_service
# - governance_service
# - identity_service
# - notifications_service
# - object_storage_service
# - profile
# - projects_service
# - source_code_hosting

- name: Wipe one managed PostgreSQL database
  hosts: all
  become: true
  force_handlers: true
  vars:
    wipe_pg_db_supported:
      - billing
      - sandbox_rental
      - mailbox_service
      - governance_service
      - identity_service
      - notifications_service
      - object_storage_service
      - profile
      - projects_service
      - source_code_hosting
  pre_tasks:
    - name: Assert supported PostgreSQL wipe target
      ansible.builtin.assert:
        that:
          - (wipe_pg_db_name | default('')) in wipe_pg_db_supported
        fail_msg: >-
          wipe_pg_db_name must be one of {{ wipe_pg_db_supported | join(', ') }}.

    - name: Check deployed PostgreSQL client binary for wipe
      ansible.builtin.stat:
        path: "{{ verself_bin }}/psql"
      register: wipe_pg_db_psql_binary

    - name: Assert deployed PostgreSQL client is present for wipe
      ansible.builtin.assert:
        that:
          - wipe_pg_db_psql_binary.stat.exists
        fail_msg: >-
          wipe-pg-db runs against an existing deployment, but {{ verself_bin }}/psql is missing.
          Run a deploy playbook before wiping a managed PostgreSQL database.

    - name: Wait for deployed PostgreSQL before wipe
      ansible.builtin.wait_for:
        host: "{{ topology_endpoints.postgresql.host }}"
        port: "{{ topology_endpoints.postgresql.endpoints.postgres.port }}"
        timeout: 5

    - name: Verify PostgreSQL admin connectivity before wipe
      ansible.builtin.command:
        argv:
          - "{{ verself_bin }}/psql"
          - -h
          - "{{ topology_endpoints.postgresql.host }}"
          - -p
          - "{{ topology_endpoints.postgresql.endpoints.postgres.port | string }}"
          - -U
          - postgres
          - -d
          - postgres
          - -v
          - ON_ERROR_STOP=1
          - -tAc
          - SELECT 1
      environment:
        PGPASSWORD: "{{ postgresql_admin_password }}"
      changed_when: false
      no_log: true

  tasks:
    - name: Resolve component name for PostgreSQL wipe target
      ansible.builtin.set_fact:
        wipe_pg_db_component_name: "{{ wipe_pg_db_component_names[wipe_pg_db_name] }}"
      vars:
        wipe_pg_db_component_names:
          billing: billing
          sandbox_rental: sandbox_rental
          mailbox_service: mailbox_service
          governance_service: governance_service
          identity_service: identity_service
          notifications_service: notifications_service
          object_storage_service: object_storage_service
          profile: profile_service
          projects_service: projects_service
          source_code_hosting: source_code_hosting_service

    - name: Select component for PostgreSQL wipe target
      ansible.builtin.set_fact:
        wipe_pg_db_component: >-
          {{
            topology_components
            | selectattr('name', 'equalto', wipe_pg_db_component_name)
            | list
            | first
          }}

    - name: Resolve Nomad jobs to stop for PostgreSQL wipe target
      ansible.builtin.set_fact:
        wipe_pg_db_stop_nomad_jobs: "{{ wipe_pg_db_stop_nomad_job_map[wipe_pg_db_name] }}"
      vars:
        wipe_pg_db_stop_nomad_job_map:
          billing:
            - sandbox-rental
            - billing
          sandbox_rental:
            - sandbox-rental
          mailbox_service:
            - mailbox-service
          governance_service:
            - governance-service
          identity_service:
            - identity-service
          notifications_service:
            - notifications-service
            - electric-notifications
          object_storage_service:
            - object-storage-admin
            - object-storage-service
          profile:
            - profile-service
          projects_service:
            - projects-service
          source_code_hosting:
            - source-code-hosting-service

    - name: Resolve substrate systemd units to stop for PostgreSQL wipe target
      ansible.builtin.set_fact:
        wipe_pg_db_stop_systemd_units: "{{ wipe_pg_db_stop_systemd_unit_map[wipe_pg_db_name] }}"
      vars:
        wipe_pg_db_stop_systemd_unit_map:
          billing: []
          sandbox_rental:
            - electric
          mailbox_service: []
          governance_service: []
          identity_service: []
          notifications_service:
            - electric-notifications
          object_storage_service: []
          profile: []
          projects_service: []
          source_code_hosting: []

    - name: Probe Nomad jobs before PostgreSQL wipe target
      ansible.builtin.command:
        argv:
          - "{{ verself_bin }}/nomad"
          - job
          - status
          - "{{ item }}"
      environment:
        NOMAD_ADDR: http://127.0.0.1:4646
      changed_when: false
      failed_when: false
      loop: "{{ wipe_pg_db_stop_nomad_jobs }}"
      loop_control:
        label: "{{ item }}"
      register: wipe_pg_db_nomad_job_status

    - name: Stop Nomad jobs before PostgreSQL wipe target
      ansible.builtin.command:
        argv:
          - "{{ verself_bin }}/nomad"
          - job
          - stop
          - -detach
          - "{{ item.item }}"
      environment:
        NOMAD_ADDR: http://127.0.0.1:4646
      loop: "{{ wipe_pg_db_nomad_job_status.results | selectattr('rc', 'equalto', 0) | list }}"
      loop_control:
        label: "{{ item.item }}"
      register: wipe_pg_db_nomad_job_stop
      changed_when: true

    - name: Stop substrate systemd units before PostgreSQL wipe target
      ansible.builtin.systemd:
        name: "{{ item }}"
        state: stopped
      loop: "{{ wipe_pg_db_stop_systemd_units }}"
      loop_control:
        label: "{{ item }}"
      failed_when: false

    - name: Include nftables handler role for component wipe
      ansible.builtin.include_role:
        name: nftables

    - name: Classify PostgreSQL wipe target ClickHouse binding
      ansible.builtin.set_fact:
        wipe_pg_db_uses_clickhouse: >-
          {{ wipe_pg_db_component.workload.clickhouse is defined }}

    - name: Load deployed ClickHouse operator defaults for PostgreSQL wipe target
      ansible.builtin.include_vars:
        file: ../roles/clickhouse/defaults/main.yml
      when: wipe_pg_db_uses_clickhouse | bool

    - name: Check deployed ClickHouse client binary for PostgreSQL wipe target
      ansible.builtin.stat:
        path: "{{ verself_bin }}/clickhouse-client"
      register: wipe_pg_db_clickhouse_client_binary
      when: wipe_pg_db_uses_clickhouse | bool

    - name: Assert deployed ClickHouse client is present for PostgreSQL wipe target
      ansible.builtin.assert:
        that:
          - wipe_pg_db_clickhouse_client_binary.stat.exists
        fail_msg: >-
          wipe-pg-db reconciles ClickHouse grants for {{ wipe_pg_db_name }},
          but {{ verself_bin }}/clickhouse-client is missing. Run aspect deploy first.
      when: wipe_pg_db_uses_clickhouse | bool

    - name: Set deployed ClickHouse operator client command for PostgreSQL wipe target
      ansible.builtin.set_fact:
        clickhouse_operator_client_command: >-
          sudo -u {{ clickhouse_operator_system_user }} {{ verself_bin }}/clickhouse-client
          --config-file {{ clickhouse_operator_client_config_path }}
          --user {{ clickhouse_operator_db_user }}
      when: wipe_pg_db_uses_clickhouse | bool

    - name: Wait for deployed ClickHouse before PostgreSQL wipe target reconciliation
      ansible.builtin.wait_for:
        host: "{{ topology_endpoints.clickhouse.host }}"
        port: "{{ topology_endpoints.clickhouse.endpoints.native_tls.port }}"
        timeout: 5
      when: wipe_pg_db_uses_clickhouse | bool

    - name: Verify deployed ClickHouse operator access before PostgreSQL wipe target reconciliation
      ansible.builtin.command:
        cmd: "{{ clickhouse_operator_client_command }} --query 'SELECT 1'"
      changed_when: false
      when: wipe_pg_db_uses_clickhouse | bool

    - name: Include TigerBeetle for billing PostgreSQL wipe
      ansible.builtin.include_role:
        name: tigerbeetle
      when: wipe_pg_db_name == 'billing'

    - name: Reset selected PostgreSQL database
      ansible.builtin.include_role:
        name: postgresql_db_reset
      vars:
        postgresql_db_reset_name: "{{ wipe_pg_db_component.postgres.database }}"
        postgresql_db_reset_owner: "{{ wipe_pg_db_component.postgres.owner }}"

    - name: Reset TigerBeetle for billing PostgreSQL wipe
      when: wipe_pg_db_name == 'billing'
      block:
        - name: Stop TigerBeetle before billing ledger wipe
          ansible.builtin.systemd:
            name: tigerbeetle
            state: stopped

        - name: Remove TigerBeetle data file before billing ledger wipe
          ansible.builtin.file:
            path: /var/lib/tigerbeetle/data.tigerbeetle
            state: absent

        - name: Format TigerBeetle data file after billing ledger wipe
          ansible.builtin.command:
            cmd: >-
              {{ verself_bin }}/tigerbeetle format
              --cluster=0
              --replica=0
              --replica-count=1
              /var/lib/tigerbeetle/data.tigerbeetle
          become: true
          become_user: tigerbeetle
          changed_when: true

        - name: Start TigerBeetle after billing ledger wipe
          ansible.builtin.systemd:
            name: tigerbeetle
            state: started

        - name: Wait for TigerBeetle after billing ledger wipe
          ansible.builtin.wait_for:
            host: "{{ topology_endpoints.tigerbeetle.host }}"
            port: "{{ topology_endpoints.tigerbeetle.endpoints.client.port }}"
            timeout: 30

    - name: Clear Electric shape storage after sandbox_rental PostgreSQL wipe
      ansible.builtin.file:
        path: /var/lib/electric
        state: absent
      when: wipe_pg_db_name == 'sandbox_rental'

    - name: Recreate empty Electric shape storage after sandbox_rental PostgreSQL wipe
      ansible.builtin.file:
        path: /var/lib/electric
        state: directory
        owner: root
        group: root
        mode: "0750"
      when: wipe_pg_db_name == 'sandbox_rental'

    - name: Configure selected service substrate after PostgreSQL wipe
      ansible.builtin.include_tasks:
        file: ../tasks/component_substrate.yml
      vars:
        component: "{{ wipe_pg_db_component }}"

    - name: Include Electric after sandbox_rental PostgreSQL wipe
      ansible.builtin.include_role:
        name: electric
      vars:
        electric_config: "{{ topology_electric_instances.default }}"
      when: wipe_pg_db_name == 'sandbox_rental'

    - name: Include Electric Notifications after notifications_service PostgreSQL wipe
      ansible.builtin.include_role:
        name: electric
      vars:
        electric_config: "{{ topology_electric_instances.notifications }}"
      when: wipe_pg_db_name == 'notifications_service'
