# src/substrate/ansible/roles/firecracker/templates/vm-orchestrator-seed.service.j2
[Unit]
Description=Seed Firecracker composable image zvols via vm-orchestrator
# Runs once after the daemon socket is up. Each ExecStart= is a single
# SeedImage RPC; the daemon owns staging zvol, dd/mkfs, snapshot, and
# dependent-clone teardown. ExecStart= entries are emitted in tier order (substrate →
# platform_toolchain → customer_uploaded) so dependents land after
# their bases — Jinja iterates the authored list directly because
# config.cue declares the catalog in that order.
After=vm-orchestrator.service
Requires=vm-orchestrator.service
PartOf=vm-orchestrator.service

[Service]
Type=oneshot
RemainAfterExit=yes
{% for image in firecracker_seed_images %}
# tier={{ image.tier }}
ExecStart={{ firecracker_vm_orchestrator_cli_path }} seed-image \
  --socket {{ firecracker_vm_orchestrator_api_socket }} \
  --ref {{ image.ref }} \
  --strategy {{ image.strategy }} \
  --size-bytes {{ image.size_bytes }} \
  --volblocksize {{ image.volblocksize }} \
{% if image.source_path | length > 0 %}
  --source-path {{ image.source_path }} \
{% endif %}
{% if image.filesystem_label | length > 0 %}
  --filesystem-label {{ image.filesystem_label }} \
{% endif %}
  --allow-destroying-active-clones
{% endfor %}
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
NoNewPrivileges=yes

[Install]
WantedBy=multi-user.target
