# src/platform/ansible/roles/nomad/templates/nomad.service.j2
# {{ ansible_managed }}
[Unit]
Description=Nomad
Documentation=https://developer.hashicorp.com/nomad/docs
After=verself-firewall.target network.target
Requires=verself-firewall.target

[Service]
Type=notify
User={{ nomad_user }}
Group={{ nomad_group }}
Environment=OTEL_SERVICE_NAME=nomad
ExecReload=/bin/kill -HUP $MAINPID
ExecStart={{ verself_bin }}/nomad agent -config={{ nomad_config_dir }}/nomad.hcl
KillMode=process
KillSignal=SIGINT
LimitNOFILE=65536
LimitNPROC=infinity
TasksMax=infinity
OOMScoreAdjust=-1000
Restart=on-failure
RestartSec=2
# Delegate the cgroup-v2 hierarchy to Nomad so it can create child
# cgroups for each alloc. Without this, raw_exec tasks running as a
# non-root user (e.g. profile_service) fail at fork/exec with
# "permission denied" — systemd's cgroup ownership keeps the alloc
# scope unwritable by the unprivileged child even though Nomad itself
# runs as root.
Delegate=yes

[Install]
WantedBy=multi-user.target
