# src/platform/ansible/roles/platform/templates/platform.service.j2
[Unit]
Description=Forge Metal Platform Docs Frontend (TanStack Start)
After=forge-metal-firewall.target network.target
Requires=forge-metal-firewall.target

[Service]
Type=simple
User={{ platform_runtime_user }}
Group={{ platform_runtime_group }}
WorkingDirectory={{ platform_app_dir }}
EnvironmentFile=/etc/platform/env
ExecStart={{ fm_bin }}/node --import ./instrumentation.mts .output/server/index.mjs
Restart=on-failure
RestartSec=5
Environment=PATH={{ fm_bin }}:/usr/bin:/bin
Environment=HOME={{ platform_state_dir }}
CapabilityBoundingSet=
ProtectHome=true
ProtectSystem=strict
ReadWritePaths={{ platform_state_dir }}
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
LockPersonality=true
NoNewPrivileges=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=multi-user.target
