# src/platform/ansible/roles/platform/templates/platform.service.j2
[Unit]
Description=Forge Metal Platform Docs Frontend (TanStack Start)
After=forge-metal-firewall.target network.target
Requires=forge-metal-firewall.target

[Service]
Type=simple
User={{ ansible_user }}
WorkingDirectory={{ platform_app_dir }}
EnvironmentFile=/etc/platform/env
ExecStart={{ fm_bin }}/node --import ./instrumentation.mts .output/server/index.mjs
Restart=on-failure
RestartSec=5
Environment=PATH={{ fm_bin }}:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
