# src/platform/ansible/roles/webmail/templates/webmail.service.j2
[Unit]
Description=Webmail Frontend (TanStack Start)
After=forge-metal-firewall.target network.target mailbox-service.service electric-mail.service
Requires=forge-metal-firewall.target
Wants=mailbox-service.service electric-mail.service

[Service]
Type=simple
User={{ webmail_runtime_user }}
Group={{ webmail_runtime_group }}
WorkingDirectory={{ webmail_app_dir }}
EnvironmentFile=/etc/webmail/env
ExecStart={{ fm_bin }}/node --import ./instrumentation.mts .output/server/index.mjs
Restart=on-failure
RestartSec=5
Environment=PATH={{ fm_bin }}:/usr/bin:/bin
Environment=HOME={{ webmail_state_dir }}
CapabilityBoundingSet=
ProtectHome=true
ProtectSystem=strict
ReadWritePaths={{ webmail_state_dir }}
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
LockPersonality=true
NoNewPrivileges=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=multi-user.target
