# src/components/bazel_remote/templates/bazel-remote.service.j2
[Unit]
Description=bazel-remote build cache
After=verself-firewall.target network.target zfs-mount.service
Requires=verself-firewall.target

[Service]
Type=simple
User={{ bazel_remote_user }}
Group={{ bazel_remote_group }}
Environment=OTEL_SERVICE_NAME=bazel-remote
WorkingDirectory={{ bazel_remote_state_dir }}
# bazel-remote listens on every interface so wg-ops peers can reach it; the
# host firewall enforces that only lo + wg-ops are accepted at L3.
ExecStart={{ verself_bin }}/bazel-remote \
    --dir {{ bazel_remote_state_dir }} \
    --max_size {{ bazel_remote_max_size_gib }} \
    --max_size_hard_limit {{ bazel_remote_max_size_hard_limit_gib }} \
    --grpc_address 0.0.0.0:{{ bazel_remote_grpc_port }} \
    --http_address 0.0.0.0:{{ bazel_remote_http_port }} \
    --enable_endpoint_metrics \
    --http_metrics_prefix
Restart=on-failure
RestartSec=5
LimitNOFILE=65536

# Hardening — note: RestrictSUIDSGID=true is intentionally NOT set. bazel-remote
# opens cache files with O_CREAT and a mode whose setgid bit is set (mode
# 02664) so newly-written blobs inherit the cache's group ownership; with
# RestrictSUIDSGID=true the kernel returns EPERM and every cache write fails
# silently with INTERNAL gRPC errors. Cache files contain build artefacts and
# action results, never executables run by another user, so the suid/sgid
# vector this option mitigates does not apply here.
CapabilityBoundingSet=
NoNewPrivileges=true
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectHome=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectSystem=strict
LockPersonality=true
ReadWritePaths={{ bazel_remote_state_dir }}
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=true
RestrictRealtime=true
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=multi-user.target
