# src/platform/ansible/roles/nats/templates/nats.service.j2
[Unit]
Description=NATS JetStream
After=forge-metal-firewall.target network.target spire-agent.service nats-spiffe-helper.service
Requires=forge-metal-firewall.target
Wants=spire-agent.service nats-spiffe-helper.service

[Service]
Type=simple
User=nats
Group=nats
SupplementaryGroups={{ spire_workload_group }}
ExecStart={{ fm_bin }}/nats-server -c {{ nats_config_dir }}/nats-server.conf
ExecReload={{ fm_bin }}/nats-server --signal reload=/run/nats/nats.pid
Restart=on-failure
RestartSec=5
LimitNOFILE=131072
RuntimeDirectory=nats
RuntimeDirectoryMode=0750
ProtectHome=true
ProtectSystem=strict
ReadWritePaths={{ nats_data_dir }} /run/nats
PrivateDevices=true
PrivateTmp=true
ProtectClock=true
ProtectControlGroups=true
ProtectKernelLogs=true
ProtectKernelModules=true
ProtectKernelTunables=true
LockPersonality=true
NoNewPrivileges=true
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native
UMask=0077

[Install]
WantedBy=multi-user.target
